# sysroot/etc/systemd/system/lsgp@.service
[Unit]
Description=LSGP virtual machine starter

[Service]
ExecStart=/usr/bin/bash -c 'lsgp %i'

[Install]
WantedBy=multi-user.target
